# services/systemd/scx.timer
[Unit]
Description=Run scx schedulers

[Timer]
OnBootSec=30sec
Unit=scx.service

[Install]
WantedBy=graphical.target
